# upstart/fdfs_storaged.conf
description "fastdfs storage daemon"

start on runlevel [23]
stop on shutdown
env DAEMON=false
respawn

exec /usr/bin/fdfs_storaged /etc/fdfs/storage.conf

