# upstart/fdfs_trackerd.conf
description "fastdfs tracker daemon"

start on runlevel [23]
stop on shutdown
env DAEMON=false
respawn

exec /usr/bin/fdfs_trackerd /etc/fdfs/tracker.conf

